# etc/systemd/system/numlock-on-boot.service
[Unit]
Description=Activate numlock on system boot

[Service]
ExecStart=/usr/local/bin/numlock
StandardInput=tty
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

